# controlplane/pipecat-voice/systemd/nova-tui-bridge@.service
[Unit]
Description=Nova TUI NATS Bridge for %i
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/adapt/platform/novaops/controlplane/pipecat-voice
EnvironmentFile=-/adapt/secrets/db.env
EnvironmentFile=-/adapt/secrets/m2.env
Environment=DISPLAY=:0
Environment=SUBJECT_NS=nova
Environment=BRIDGE_RUNTIME=tui
Environment=BRIDGE_LEGACY_SUBJECTS=1
Environment=ECHO_TUI_AGENT=%i
Environment=ECHO_TUI_WINDOW_NAME=%i CLI
Environment=ECHO_PROFILE_ROOT=/home/x/.hermes/profiles/%i
Environment=ECHO_TUI_ALLOW_CLI_FALLBACK=1
Environment=ECHO_TUI_DELIVERY_TIMEOUT=20
Environment=ECHO_TUI_REPLY_CAPTURE_TIMEOUT=45
Environment=ECHO_TUI_MAX_QUEUE_DEPTH=32
ExecStart=/usr/bin/python3 /adapt/platform/novaops/controlplane/pipecat-voice/scripts/echo_tui_nats_bridge.py
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
