# controlplane/pipecat-voice/systemd/nova-crew-route-state.service
[Unit]
Description=Nova crew route-state snapshot
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/adapt/platform/novaops/controlplane/pipecat-voice
Environment=DISPLAY=:0
ExecStart=/usr/bin/python3 /adapt/platform/novaops/controlplane/pipecat-voice/scripts/crew_route_state.py
